# phase3-monitoring/dcgm_exporter.service
[Unit]
Description=NVIDIA DCGM Exporter
