# phase3-monitoring/node_exporter.service
[Unit]
Description=Prometheus Node Exporter
