# baseline_logs/2025-05-18/baseline.timer
[Unit]
Description=Monthly System Baseline

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target